# templates/upstart.conf
start on [345]
stop on [06]

respawn
chdir {{ project_dir }}
exec supervisord -u {{ username }} -c {{ supervisord_conf }} -n 2>> {{ project_dir }}/logs/supervisord.error.log
